# modules/vulnerabilities/unix/http/apache_spark_rce/templates/spark-master.service.erb
[Unit]
Description=Apache Spark Master
After=network.target

[Service]
Type=forking
User=<%= @user %>
Environment="PATH=/usr/local/bin:/usr/bin:/bin:/usr/local/games:/usr/games:/usr/local/spark/sbin"
Environment="JAVA_HOME=/usr/lib/jvm/java-11-openjdk"
WorkingDirectory=/usr/local/spark/bin/
ExecStart=/usr/local/spark/sbin/start-master.sh --host 0.0.0.0 --properties-file /usr/local/spark/conf/spark-defaults.conf
ExecStop=/usr/local/spark/sbin/stop-master.sh
Restart=on-abort
RestartSec=1

[Install]
WantedBy=multi-user.target