# modules/vulnerabilities/unix/http/apache_http_server_2_4_4_9_path_traversal_and_remote_code_execution/files/httpd.service
# Custom Apache Web Server service file.
[Unit]
Description=Apache HTTPD 2.4.49 Web Server
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=/usr/local/apache2/bin/httpd
WorkingDirectory=/usr/local/apache2/bin/

[Install]
WantedBy=multi-user.target